# scripts/etc/crontabs/root
# 定时重启（每个星期六凌晨2点50）
#50 2 * * 6 sleep 5 && touch /etc/banner && reboot
# 定时重启网络-更换IP
50 2 * * 6 /etc/init.d/network restart
52 2 * * 6 /etc/init.d/zerotier restart
